# services/runit/finish
#!/bin/sh
# runit finish script for anvillm
# Called when the service stops

echo "anvillm stopped"

# Remove PID file if it exists
rm -f /tmp/anvillm.pid

# Optional: clean up stale sockets
# export NAMESPACE="${NAMESPACE:-/tmp/ns.$(whoami)}"
# rm -f "${NAMESPACE}/agent"
